# rootfs/etc/s6-overlay/s6-rc.d/initialization/run
#!/command/with-contend bashio

# CUPS paths
readonly cups_config_path=/config/cups
readonly src_cups_templates_path=/usr/templates/cups
readonly cups_log_path=/config/cups/logs
readonly cups_client_cfg=client.conf.tempio
readonly cups_daemon_cfg=cupsd.conf.tempio
readonly cups_files_cfg=cups-files.conf.tempio

# Avahi paths
readonly avahi_config_path=/config/avahi
readonly avahi_templates_path=/config/templates/avahi
readonly cups_templates_path=/config/templates/cups
readonly src_avahi_templates_path=/usr/templates/avahi
readonly avahi_daemon_cfg=avahi-daemon.conf.tempio

mkdir -p /var/run/dbus

# Install user configured/requested packages
if bashio::config.has_value 'packages'; then
  apt update ||
    bashio::exit.nok 'Failed updating packages repository indexes'

  for package in $(bashio::config 'packages'); do
    apt-get install -y "$package" ||
      bashio::exit.nok "Failed installing package ${package}"
  done
fi

# Cups config folder
if ! bashio::fs.directory_exists "${cups_config_path}"; then
  mkdir -p "${cups_config_path}" ||
    bashio::exit.nok 'Failed to create a persistent cups config folder'

  chmod 700 "${cups_config_path}" ||
    bashio::exit.nok \
      'Failed setting permissions on persistent cups config folder'
fi

# Cups log folder
if ! bashio::fs.directory_exists "${cups_log_path}"; then
  mkdir -p "${cups_log_path}" ||
    bashio::exit.nok 'Failed to create a persistent cups logging folder'

  chmod 700 "${cups_log_path}" ||
    bashio::exit.nok \
      'Failed setting permissions on persistent cups logging folder'
fi

# Avahi config folder
if ! bashio::fs.directory_exists "${avahi_config_path}"; then
  mkdir -p "${avahi_config_path}" ||
    bashio::exit.nok 'Failed to create a persistent Avahi config folder'

  chmod 700 "${avahi_config_path}" ||
    bashio::exit.nok \
      'Failed setting permissions on persistent Avahi config folder'
fi

# avahi templates folder
if ! bashio::fs.directory_exists "${avahi_templates_path}"; then
  mkdir -p "${avahi_templates_path}" ||
    bashio::exit.nok 'Failed to create a persistent avahi templates folder'

  chmod 700 "${avahi_templates_path}" ||
    bashio::exit.nok \
      'Failed setting permissions on persistent avahi templates folder'
fi

# cups templates folder
if ! bashio::fs.directory_exists "${cups_templates_path}"; then
  mkdir -p "${cups_templates_path}" ||
    bashio::exit.nok 'Failed to create a persistent cups templates folder'

  chmod 700 "${cups_templates_path}" ||
    bashio::exit.nok \
      'Failed setting permissions on persistent cups templates folder'
fi

if [ ! -e "$avahi_templates_path/$avahi_daemon_cfg".tempio ]; then
  cp "$src_avahi_templates_path/$avahi_daemon_cfg" "$avahi_templates_path/$avahi_daemon_cfg"
fi

if [ ! -e "$cups_templates_path/$cups_client_cfg".tempio ]; then
  cp "$src_cups_templates_path/$cups_client_cfg" "$cups_templates_path/$cups_client_cfg"
fi

if [ ! -e "$cups_templates_path/$cups_daemon_cfg".tempio ]; then
  cp "$src_cups_templates_path/$cups_daemon_cfg" "$cups_templates_path/$cups_daemon_cfg"
fi

if [ ! -e "$cups_templates_path/$cups_files_cfg".tempio ]; then
  cp "$src_cups_templates_path/$cups_files_cfg" "$cups_templates_path/$cups_files_cfg"
fi
