# rootfs/etc/s6-overlay/s6-rc.d/initialization/up
#!/command/with-contend bashio


mkdir -p /var/run/dbus
