# rootfs/etc/s6-overlay/s6-rc.d/init-config/run
#!/command/with-contend bashio
# shellcheck shell=bash

# ==============================================================================
# Home Assistant  Add-on: Cups Server
# Persists user settings and installs custom user packages.
# ==============================================================================
# local file paths
templates_path="/templates"
cups_cfg_path="/etc/cups/"
cups_ssl_path="$cups_cfg_path/ssl"
avahi_cfg_path="/etc/avahi"


# HA File paths
real_config_path="/data"
real_ssl_path="/ssl"

# CUPS paths
real_cups_path="$real_config_path/cups"
cups_client_cfg="client.conf"
cups_daemon_cfg="cupsd.conf"
cups_files_cfg="cups-files.conf"

# Avahi paths
real_avahi_path="$real_config_path/avahi"
avahi_daemon_cfg="avahi-daemon.conf"

if [ ! -d "$real_config_path" ]; then
  mkdir -p "$real_config_path"
fi

# CUPS Files
if [ ! -d "$real_cups_path" ]; then
  mkdir -p "$real_cups_path"
fi

if [ -L "$cups_cfg_path" ] && [ -d "$cups_cfg_path" ]; then
  echo "$cups_cfg_path is a valid symbolic link"
else
  rm "$cups_cfg_path"
  ln -s "$cups_cfg_path" "$real_cups_path"
fi

if [ ! -e "$cups_cfg_path/$cups_client_cfg" ]; then
  cp "$templates_path/$cups_client_cfg" "$cups_cfg_path/$cups_client_cfg"
fi

if [ ! -e "$cups_cfg_path/$cups_daemon_cfg" ]; then
  cp "$templates_path/$cups_daemon_cfg" "$cups_cfg_path/$cups_daemon_cfg"
fi

if [ ! -e "$cups_cfg_path/$cups_files_cfg" ]; then
  cp "$templates_path/$cups_files_cfg" "$cups_cfg_path/$cups_files_cfg"
fi

if [ -L "$cups_ssl_path" ] && [ -d "$cups_ssl_path" ]; then
  echo "$cups_ssl_path is a valid symbolic link"
else
  rm "$cups_ssl_path"
  ln -s "$cups_ssl_path" "$real_ssl_path"
fi

# AVAHI Files
if [ ! -d "$real_avahi_path" ]; then
  mkdir -p "$real_avahi_path"
fi

if [ -L "$avahi_cfg_path" ] && [ -d "$avahi_cfg_path" ]; then
  echo "$avahi_cfg_path is a valid symbolic link"
else
  rm "$avahi_cfg_path"
  ln -s "$avahi_cfg_path" "$real_avahi_path"
fi

if [ ! -e "$avahi_cfg_path/$avahi_daemon_cfg" ]; then
  cp "$templates_path/$avahi_daemon_cfg" "$avahi_cfg_path/$avahi_daemon_cfg"
fi
